# systemd/user/gwe.service
[Unit]
Description=Green with envy

[Service]
Type=simple
ExecStart=/usr/bin/gwe
Restart=on-failure

[Install]
WantedBy=default.target